# root/etc/services.d/aria2c/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc /app/forego start -f /app/Procfile
